# root/etc/services.d/cloud-torrent/run
#!/usr/bin/with-contenv bash

chmod a+x /defaults/cloud-torrent

exec \
	s6-setuidgid abc /defaults/cloud-torrent -p $PORT