# root/etc/services.d/cloud-torrent/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid abc /defaults/cloud-torrent --port $PORT