# root/etc/services.d/cloud-torrent/run
#!/usr/bin/with-contenv bash

chmod a+x /usr/local/bin/cloud-torrent

exec \
	s6-setuidgid abc /usr/local/bin/cloud-torrent -p $PORT -a $AUTH -c /config/cloud-torrent.json